# deploy/blackcat.service
[Unit]
Description=BlackCat AI Agent Daemon
After=network-online.target

[Service]
Type=simple
ExecStart=__BLACKCAT_BINARY__ daemon --config __DEPLOY_CONFIG_PATH__
WorkingDirectory=__DEPLOY_HOME__
Restart=on-failure
RestartSec=5
Environment=BLACKCAT_VAULT_PASSPHRASE=__VAULT_PASSPHRASE__

[Install]
WantedBy=default.target
