# deploy/opencode.service
[Unit]
Description=OpenCode Server (managed by BlackCat)
After=network-online.target

[Service]
Type=simple
ExecStart=__BLACKCAT_BINARY__ serve --port __OPENCODE_PORT__ --binary __OPENCODE_BINARY__
WorkingDirectory=__DEPLOY_HOME__
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
